# systemd/telegram-bridge-watcher.path
[Unit]
Description=Watch bridge.py and hooks for changes (hot reload)

[Path]
PathModified=%BRIDGE_DIR%/bridge.py
PathModified=%BRIDGE_DIR%/hooks

[Install]
WantedBy=default.target
